# crontab
MAILTO="root"

30 0 * * * /srv/applewiki/cron-daily.sh
